# discord-bot.service
[Unit]
Description=Discord Content Announcement Bot Service
After=network.target

[Service]
Type=simple
User=%i
Environment="DISPLAY=:99"
ExecStart=%h/discord-bot/scripts/start-bot.sh
Restart=on-failure
RestartSec=10s
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
